# config/autostart/conky.desktop
[Desktop Entry]
Comment=Conky
Exec=conky
GenericName=Conky
Hidden=true
Name=Conky
StartupNotify=false
Terminal=false
Type=Application
