# config/autostart/conky.desktop
[Desktop Entry]
Name=Conky
GenericName=Conky
Comment=Conky
Exec=conky
Terminal=false
Type=Application
StartupNotify=false
